# layers/linux/HOME/.config/systemd/user/clipster.service
[Unit]
Description=clipster daemon

[Service]
Type=simple
ExecStart=/usr/bin/clipster --daemon
Restart=always

[Install]
WantedBy=default.target